# shells/ebaytoold
#!/bin/sh
#
# ebaytoold
#
# chkconfig: - 85 15
# description: ebaytool daemon
#

# todo: which is better? pool18 count.
CPBASE=/var/www/ebaytool.jp/build/WEB-INF
#CPBASE=/usr/local/apache-tomcat/webapps/ROOT/WEB-INF

PATH=/usr/local/jdk/bin:$PATH
export CLASSPATH=$CPBASE/classes:$CPBASE/lib/*

start() {
	java ebaytool.apicall.Daemon start >> /var/www/ebaytool.jp/logs/daemon.log &
}

stop() {
	java ebaytool.apicall.Daemon stop
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		start
		;;
	*)
		echo $"Usage: $prog {start|stop|restart}"
esac
