# shells/ebaytoold
#!/bin/sh
#
# ebaytoold
#
# chkconfig: - 85 15
# description: ebaytool daemon
#

BASEDIR="$( dirname "$( cd "$( dirname "${BASH_SOURCE[0]}" )" && pwd )" )"
BASENAME=${BASEDIR/*\//}

cd $BASEDIR

CPBASE=$BASEDIR/build/WEB-INF

PATH=/usr/local/jdk/bin:$PATH
export CLASSPATH=$CPBASE/classes:$CPBASE/lib/*

start() {
	# option "-Xms1000m" is needed to read large xml files.
	#java -Xms2500m ebaytool.apicall.Daemon start >> $BASEDIR/logs/daemon.log 2>&1 &
	java -Xmx1000m ebaytool.apicall.Daemon start >> $BASEDIR/logs/daemon.log 2>&1 &
	#java ebaytool.apicall.Daemon start >> $BASEDIR/logs/daemon.log 2>&1 &
}

stop() {
	java ebaytool.apicall.Daemon stop
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		start
		;;
	*)
		echo $"Usage: $prog {start|stop|restart}"
esac
